# sheepit-binary-docker/crontab
# Calls "sync-mirror" (aka sync.sh) every 6 hours
0 */6 * * * sync-mirror > /dev/stdout